# etc/service/nginx/run
#! /bin/sh

/usr/sbin/nginx

while true; do sleep 60; done