# etc/service/php5-fpm/run
#! /bin/sh

/usr/sbin/php5-fpm

while true; do sleep 60; done