# etc/service/php5-fpm/run
#! /bin/sh

/usr/sbin/php5-fpm --nodaemonize
